# trash/gnome-sound-applet.service
[Service]
ExecStart=/usr/bin/gnome-sound-applet
PartOf=openbox.service

[Install]
WantedBy=default.target
