# trash/openbox.target
[Unit]
Description=Openbox Desktop Environment
Wants=xorg.target
Wants=openbox.service
AllowIsolate=true

[Install]
WantedBy=default.target
