# openbox.target
[Unit]
Description=Openbox Desktop Environment
Wants=xorg.target
Requires=dbus.socket
AllowIsolate=true

[Install]
WantedBy=default.target
